# .config/systemd/user/irb-vpn-reminder.service
[Unit]
Description=IRB VPN Reminder

[Service]
ExecStart=%h/scripts/vpn-reminder IRB
