# .config/systemd/user/irb-vpn-reminder.timer
[Unit]
Description=IRB VPN Reminder Timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
