# .config/systemd/user/uv-cache-prune.service
[Unit]
Description=Prune uv cache

[Service]
Type=oneshot
ExecStart=uv cache prune
