# .config/systemd/user/uv-cache-prune.timer
[Unit]
Description=Run uv cache prune every 4 weeks

[Timer]
OnCalendar=monthly
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
